# vapour-os-xfce/files/etc/profile.d/0000-vapour-os-xfce.sh
#!/bin/sh
DISTRO_ID="vapour-os"
cp -nr /usr/share/$DISTRO_ID/custom-configs/$DISTRO_ID-xfce/.config /home/$(whoami)/
