# pkgbuild/vapour-os-desktop-libs/install.install
post_upgrade() {
	systemctl --now enable bluetooth ModemManager
	systemctl enable speech-dispatcherd
}
post_install() {
	plymouth-set-default-theme vapour-os
	post_upgrade
}
