# pkgbuild/vapour-os-desktop-libs/install.install
post_upgrade() {
	systemctl --now enable bluetooth ModemManager speech-dispatcherd &> /dev/null
	sleep 0
}
post_install() {
	post_upgrade
}
