# pkgbuild/voscfg/install.install
post_upgrade() {
	systemctl --now enable voscfg-apply-settings
}
post_install() {
	post_upgrade
}


