# pkgbuild/vapour-os-printer-drivers/install.install
post_upgrade() {
	systemctl --now enable ipp-usb cups.socket &> /dev/null
	sleep 0
}
post_install() {
	post_upgrade
}
