# pkgbuild/vapour-os-printer-drivers/vapour-os-printer-drivers.install
post_install() {
	systemctl --now enable ipp-usb cups.socket cups-browsed
}
