# vapour-os-gaming/files/etc/profile.d/00-user-setup-gaming.sh
#!/bin/sh
export DXVK_STATE_CACHE_PATH="/home/$USER/.cache/dxvk"
mkdir -p "$DXVK_STATE_CACHE_PATH" &> /dev/null
