# vapour-os/pkg/vapour-os/etc/profile.d/zzz-user-profile-vapour-os.sh
#!/bin/sh
[ -f "/home/$USER/.profile" ] && . "/home/$USER/.profile"
