# pkgbuild/vapour-os-kde/files/usr/share/vapour-os/custom-configs/vapour-os-kde/.config/autostart/screen-reader-info-message-vapour-os-kde.desktop
[Desktop Entry]
Type=Application
Name=Screen reader usage info
Exec=/usr/lib/vapour-os/vapour-os-kde/screen-reader-info-message.sh
Terminal=false
StartupNotify=true
