# vapour-os-gui/files/usr/lib/systemd/system/min_brightness.service
[Unit]
Description=Set brightness to minimum (0.5%) if it is too low

[Service]
Type=simple
ExecStart=/opt/vapour-os/vapour-os-gui/fix_brightness

[Install]
WantedBy=multi-user.target
