# vapour-os-gui/files/usr/lib/systemd/system/min_brightness.path
[Unit]
Description=Monitor brightness for changes

[Path]
PathChanged=/etc/vapour-os/backlight/brightness
Unit=min_brightness.service

[Install]
WantedBy=multi-user.target
